# friendly-loan.service
[Unit]
Description=Friendly Loan Web Application
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/friendly-loan
Environment=PATH=/opt/friendly-loan/venv/bin
Environment=FLASK_ENV=production
Environment=SECRET_KEY=your-secret-key-here
ExecStart=/opt/friendly-loan/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 4 --timeout 120 wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
